# supervisord/hensei-api.ini
[program:hensei-api]
command=/opt/homebrew/bin/mise run s
process_name=%(program_name)s
numprocs=1
directory=/Users/justin/Developer/Granblue/hensei-api
environment=HOME="/Users/justin",MISE_CONFIG_ROOT="/Users/justin/Developer/Granblue/hensei-api",RAILS_ENV="development"
autostart=true
autorestart=unexpected
stopsignal=TERM
user=justin
stdout_logfile=/Users/justin/Developer/Granblue/hensei-api/logs/hensei-api.stdout.log
stdout_logfile_maxbytes=500KB
stdout_logfile_backups=10
stderr_logfile=/Users/justin/Developer/Granblue/hensei-api/logs/hensei-api.stderr.log
stderr_logfile_maxbytes=500KB
stderr_logfile_backups=10
serverurl=AUTO
